# muddo_agro.service
# /etc/systemd/system/muddo_agro.service
[Unit]
Description=Muddo Agro Chemicals LTD — Flask Web Application
After=network.target

[Service]
User=muddo
Group=muddo
WorkingDirectory=/home/muddo/muddo_agro
Environment="PATH=/home/muddo/muddo_agro/venv/bin"
Environment="FLASK_ENV=production"
# Set your actual secrets here:
# Environment="MAIL_PASSWORD=your_gmail_app_password"
# Environment="GA_MEASUREMENT_ID=G-XXXXXXXXXX"
# Environment="FLASK_SECRET_KEY=change_this_to_a_long_random_string"

ExecStart=/home/muddo/muddo_agro/venv/bin/gunicorn \
    --config gunicorn.conf.py \
    app:app

# Restart on failure
Restart=on-failure
RestartSec=5s

# Logging → systemd journal (no file paths needed)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=muddo_agro

[Install]
WantedBy=multi-user.target
